# custom/04-custom/files/rpi/etc/rc.local
#!/bin/bash

if [ ! -f /etc/systemd/system/multi-user.target.wants/watchdog.service ]; then
    echo "Enabling watchdog..."
    systemctl enable watchdog.service
fi

export SCREENDIR="$HOME/.screen"

exit 0
